# deploy/passvault.service
[Unit]
Description=PassVault backend (FastAPI via gunicorn/uvicorn workers)
After=network.target

[Service]
Type=simple
User=passvault
Group=passvault
WorkingDirectory=__APP_DIR__/backend
EnvironmentFile=__APP_DIR__/.env
ExecStart=__APP_DIR__/backend/.venv/bin/gunicorn main:app \
    -k uvicorn.workers.UvicornWorker \
    -w 2 -b 127.0.0.1:8000 \
    --forwarded-allow-ips 127.0.0.1 \
    --access-logfile - --error-logfile -
Restart=on-failure
RestartSec=3

# --- Hardening ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# Only the DB/data dir is writable; everything else (code, venv, .env) is read-only.
ReadWritePaths=__APP_DIR__/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
